# calendar.service
[Unit]
Description=Calendar Server
After=network.target

[Service]
Type=simple
User=taqis
WorkingDirectory=/home/taqis/Pone/Schedule
ExecStart=/usr/bin/node /home/taqis/Pone/Schedule/server.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment="NODE_ENV=production"
Environment="PORT=3000"

[Install]
WantedBy=multi-user.target
